# udev/99-low-power.rules
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/bin/sudo -u laith /home/laith/.local/bin/battery-mode.sh"
